# packages/packages_local/hosts-updater/hosts-updater.install
#!/bin/sh
post_install() {
    systemctl daemon-reload
    systemctl enable hosts-updater.service
    systemctl start hosts-updater.service
}
