# packages_local/hosts-updater/src/hosts-updater.service
[Unit]
Description=Update hosts file

[Service]
Type=oneshot
ExecStart=/usr/bin/hosts-updater

[Install]
WantedBy=default.target
